# jupiter/az-ddns/crontab
15,45 * * * * ${PERSONAL_CLOUD_DIR}/jupiter/az-ddns/az-ddns.sh && curl -fsS -m 30 --retry 5 -o /dev/null https://${HEALTHCHECKS_URL}/azure-ddns-update
