# backup/crontab
0 0 * * * ${HOME_SERVER_DIR}/backup/backup.sh
