# jupiter/syncthing/crontab
00 22 * * * node ${PERSONAL_CLOUD_DIR}/jupiter/syncthing/syncPhoneMedia.js && curl -fsS -m 30 --retry 5 -o /dev/null https://${HEALTHCHECKS_URL}/syncthing-copy-phone-media
