# jupiter/famly-app-downloader/crontab
00 22 * * * cd ${PERSONAL_CLOUD_DIR}/jupiter/famly-app-downloader && node famly.js && curl -fsS -m 30 --retry 5 -o /dev/null https://${HEALTHCHECKS_URL}/download-nursery-photos