# jupiter/jeeb-uk/crontab
0 0 * * * cd ${WORKSPACE_DIR}/Interests/Tech/Blog && ./build-and-push.sh; [ "$?" -le 1 ] && curl -fsS -m 30 --retry 5 -o /dev/null https://${HEALTHCHECKS_URL}/build-publish-jeebuk
