# jupiter/backup/crontab
00 00 * * * ${PERSONAL_CLOUD_DIR}/jupiter/backup/backup.sh
