# jupiter/backup/crontab
0 3 * * * ${PERSONAL_CLOUD_DIR}/backup/backup.sh
